# docker/cron/crontab
# Ejecuta el scheduler cada minuto
* * * * * cd /var/www/html && php artisan schedule:run >> /dev/stdout 2>&1